# config/systemd/user/nm-applet.service
[Unit]
After=graphical-session.target

[Service]
Type=exec
ExecStart=nm-applet --indicator
Restart=on-failure
Slice=app-graphical.slice

[Install]
WantedBy=graphical-session.target
