# config/systemd/user/wl-sunset-systray.service
[Unit]
After=graphical-session.target

[Service]
Type=exec
ExecStart=%h/.config/scripts/wlsunset-systray.py
Slice=app-graphical.slice
Restart=on-failure
