# config/systemd/user/xsettingsd.service
[Unit]
After=graphical-session.target

[Service]
Type=exec
ExecStart=xsettingsd
Restart=on-failure
Slice=app-graphical.slice

[Install]
WantedBy=graphical-session.target
